# deploy/systemd/kuangbiao-kb.service
[Unit]
Description=geowiki private knowledge service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kuangbiao
Group=kuangbiao
WorkingDirectory=/opt/kuangbiao-agent
EnvironmentFile=/opt/kuangbiao-agent/.env
Environment=PYTHONPATH=/opt/kuangbiao-agent/src
ExecStart=/opt/kuangbiao-agent/.venv/bin/uvicorn mining_qa.knowledge_service:app --host 127.0.0.1 --port 18081
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/kuangbiao-agent/data

[Install]
WantedBy=multi-user.target
